# systemd/bing-wallpaper.timer
[Unit]
Description=Fetch daily Bing wallpaper

[Timer]
OnCalendar=*-*-* 00:05:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
